# source/roles/mikbotam/templates/mikbotam-polling.conf.j2
[program:mikbotam-polling]
command: php Core.php
directory: /var/www/mikbotam/Saldo
